# board/orangepi/orangepi-zero3/overlay/etc/systemd/system/my_service.service
[Unit]
Description=My Auto Start Service
After=network.target

[Service]
ExecStart=/root/MonitorSystem
WorkingDirectory=/root

[Install]
WantedBy=multi-user.target